# word-guessing-game.service
[Unit]
Description=Chinese Word Guessing Game
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/word_guessing_game/backend
Environment=HF_HUB_OFFLINE=1
Environment=TRANSFORMERS_OFFLINE=1
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/root/word_guessing_game/backend/.env
ExecStart=/root/word_guessing_game/.venv/bin/gunicorn --workers 1 --threads 2 --timeout 180 --bind 0.0.0.0:5000 app:app
Restart=on-failure
RestartSec=5
TimeoutStartSec=180

[Install]
WantedBy=multi-user.target
